# pkg/arch/edead.service
[Unit]
Description=edead web service

[Service]
WorkingDirectory=/usr/share/edead
ExecStart=/usr/bin/edead

[Install]
WantedBy=multi-user.target
